# wrench/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

; graceful timeout should always be lower than stopwaitsecs to avoid orphan gunicorn workers.
[program:{{ wrench_name }}-saashq-web]
command={{ wrench_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} --max-requests {{ gunicorn_max_requests }} --max-requests-jitter {{ gunicorn_max_requests_jitter }} -t {{ http_timeout }} --graceful-timeout 30 saashq.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/web.log
stderr_logfile={{ wrench_dir }}/logs/web.error.log
stopwaitsecs=40
killasgroup=true
user={{ user }}
directory={{ sites_dir }}
startretries={{ supervisor_startretries }}

[program:{{ wrench_name }}-saashq-schedule]
command={{ wrench_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/schedule.log
stderr_logfile={{ wrench_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ wrench_dir }}
startretries={{ supervisor_startretries }}

{% if not multi_queue_consumption %}
[program:{{ wrench_name }}-saashq-default-worker]
command={{ wrench_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/worker.log
stderr_logfile={{ wrench_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ wrench_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}
{% endif %}

[program:{{ wrench_name }}-saashq-short-worker]
command={{ wrench_cmd }} worker --queue short{{',default' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/worker.log
stderr_logfile={{ wrench_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ wrench_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}

[program:{{ wrench_name }}-saashq-long-worker]
command={{ wrench_cmd }} worker --queue long{{',default,short' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/worker.log
stderr_logfile={{ wrench_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ wrench_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}

{% for worker_name, worker_details in workers.items() %}
[program:{{ wrench_name }}-saashq-{{ worker_name }}-worker]
command={{ wrench_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/worker.log
stderr_logfile={{ wrench_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ wrench_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}
{% endfor %}


{% if not skip_redis %}
[program:{{ wrench_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/redis-cache.log
stderr_logfile={{ wrench_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}
startretries={{ supervisor_startretries }}

[program:{{ wrench_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/redis-queue.log
stderr_logfile={{ wrench_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
startretries={{ supervisor_startretries }}
{% endif %}

{% if node %}
[program:{{ wrench_name }}-node-socketio]
command={{ node }} {{ wrench_dir }}/apps/saashq/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ wrench_dir }}/logs/node-socketio.log
stderr_logfile={{ wrench_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ wrench_dir }}
startretries={{ supervisor_startretries }}
{% endif %}

[group:{{ wrench_name }}-web]
programs={{ wrench_name }}-saashq-web {%- if node -%} ,{{ wrench_name }}-node-socketio {%- endif%}


{% if multi_queue_consumption %}

[group:{{ wrench_name }}-workers]
programs={{ wrench_name }}-saashq-schedule,{{ wrench_name }}-saashq-short-worker,{{ wrench_name }}-saashq-long-worker{%- for worker_name in workers -%},{{ wrench_name }}-saashq-{{ worker_name }}-worker{%- endfor %}

{% else %}

[group:{{ wrench_name }}-workers]
programs={{ wrench_name }}-saashq-schedule,{{ wrench_name }}-saashq-default-worker,{{ wrench_name }}-saashq-short-worker,{{ wrench_name }}-saashq-long-worker{%- for worker_name in workers -%},{{ wrench_name }}-saashq-{{ worker_name }}-worker{%- endfor %}

{% endif %}

{% if not skip_redis %}
[group:{{ wrench_name }}-redis]
programs={{ wrench_name }}-redis-cache,{{ wrench_name }}-redis-queue
{% endif %}
