# wrench/config/templates/systemd/saashq-wrench-redis.target
[Unit]
After=network.target
Wants={{ wrench_name }}-redis-cache.service {{ wrench_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
